# inventory/services/sabnzbd.quadlet.d/sabnzbd.container
[Container]
HostName=sabnzbd.fearn.cloud
