# inventory/services/media-file-browser.quadlet.d/server.container
[Container]
User=1001:1001
Environment=FB_DATABASE=/database/filebrowser.db
Environment=FB_PORT=8080
